# kali_rootfs/var/lib/dpkg/info/cgmanager.conffiles
/etc/init.d/cgmanager
/etc/init.d/cgproxy
/etc/init/cgmanager.conf
/etc/init/cgproxy.conf
