# kali_rootfs/var/lib/dpkg/info/cgmanager.list
/.
/etc
/etc/init.d
/etc/init.d/cgproxy
/etc/init.d/cgmanager
/etc/init
/etc/init/cgproxy.conf
/etc/init/cgmanager.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/cgproxy.service
/lib/systemd/system/cgmanager.service
/lib/cgmanager
/lib/cgmanager/cgm-release-agent
/usr
/usr/share
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/cgm.1.gz
/usr/share/man/man8
/usr/share/man/man8/cgproxy.8.gz
/usr/share/man/man8/cgmanager.8.gz
/usr/share/doc
/usr/share/doc/cgmanager
/usr/share/doc/cgmanager/changelog.Debian.gz
/usr/share/doc/cgmanager/changelog.gz
/usr/share/doc/cgmanager/copyright
/usr/share/doc/cgmanager/README.Debian
/usr/share/upstart
/usr/share/upstart/sessions
/usr/share/upstart/sessions/notify-cgmanager.conf
/usr/bin
/usr/bin/cgm
/sbin
/sbin/cgmanager
/sbin/cgproxy
