# kali_rootfs/var/lib/dpkg/info/debtags.list
/.
/etc
/etc/cron.daily
/etc/cron.daily/debtags
/etc/apt
/etc/apt/apt.conf.d
/etc/apt/apt.conf.d/80debtags
/var
/var/lib
/var/lib/debtags
/usr
/usr/share
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/debtags
/usr/share/debtags
/usr/share/debtags/vocabulary
/usr/share/debtags/tags-current.gz
/usr/share/doc
/usr/share/doc/debtags
/usr/share/doc/debtags/paper-debtags.rst.gz
/usr/share/doc/debtags/packagesearch.png
/usr/share/doc/debtags/copyright
/usr/share/doc/debtags/README.Debian
/usr/share/doc/debtags/debtags-edit.png
/usr/share/doc/debtags/changelog.gz
/usr/share/doc/debtags/buildinfo_all.gz
/usr/share/doc/debtags/Makefile
/usr/share/doc/debtags/ranganathan.png
/usr/share/doc/debtags/README
/usr/bin
/usr/bin/debtags-submit-patch
/usr/bin/debtags-hardware
/usr/bin/debtags
