# kali_rootfs/var/lib/dpkg/info/debtags.conffiles
/etc/apt/apt.conf.d/80debtags
/etc/cron.daily/debtags
