# kali_rootfs/var/lib/dpkg/info/fake-hwclock.list
/.
/sbin
/sbin/fake-hwclock
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/fake-hwclock.service
/etc
/etc/default
/etc/default/fake-hwclock
/etc/cron.hourly
/etc/cron.hourly/fake-hwclock
/etc/init.d
/etc/init.d/fake-hwclock
/usr
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/fake-hwclock.8.gz
/usr/share/doc
/usr/share/doc/fake-hwclock
/usr/share/doc/fake-hwclock/changelog.gz
/usr/share/doc/fake-hwclock/copyright
