# kali_rootfs/var/lib/dpkg/info/john.list
/.
/usr
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/unshadow.8.gz
/usr/share/man/man8/unique.8.gz
/usr/share/man/man8/john.8.gz
/usr/share/man/man8/unafs.8.gz
/usr/share/man/man8/mailer.8.gz
/usr/share/doc
/usr/share/doc/john
/usr/share/doc/john/README.gz
/usr/share/doc/john/MODES.gz
/usr/share/doc/john/copyright
/usr/share/doc/john/RULES.gz
/usr/share/doc/john/CONFIG.gz
/usr/share/doc/john/CREDITS
/usr/share/doc/john/README.Debian
/usr/share/doc/john/changelog.gz
/usr/share/doc/john/EXAMPLES.gz
/usr/share/doc/john/FAQ.gz
/usr/share/doc/john/changelog.Debian.gz
/usr/share/doc/john/OPTIONS.gz
/usr/share/doc/john/EXTERNAL.gz
/usr/share/john
/usr/share/john/john.local.conf
/usr/share/john/john.conf
/usr/share/john/repeats32.conf
/usr/share/john/dynamic.conf
/usr/share/john/dumb16.conf
/usr/share/john/dynamic_flat_sse_formats.conf
/usr/share/john/cronjob
/usr/share/john/regex_alphabets.conf
/usr/share/john/korelogic.conf
/usr/share/john/repeats16.conf
/usr/share/john/dumb32.conf
/usr/sbin
/usr/sbin/mailer
/usr/sbin/john
/etc
/etc/cron.d
/etc/cron.d/john
/etc/john
/etc/john/john.conf
/etc/john/john-mail.conf
/etc/john/john-mail.msg
/usr/sbin/ssh2john
/usr/sbin/unshadow
/usr/sbin/keychain2john
/usr/sbin/keyring2john
/usr/sbin/racf2john
/usr/sbin/gpg2john
/usr/sbin/base64conv
/usr/sbin/pfx2john
/usr/sbin/hccap2john
/usr/sbin/dmg2john
/usr/sbin/pwsafe2john
/usr/sbin/keepass2john
/usr/sbin/unafs
/usr/sbin/putty2john
/usr/sbin/undrop
/usr/sbin/rar2john
/usr/sbin/unique
/usr/sbin/zip2john
/usr/sbin/kwallet2john
