# kali_rootfs/var/lib/dpkg/info/john.conffiles
/etc/cron.d/john
/etc/john/john-mail.conf
/etc/john/john-mail.msg
/etc/john/john.conf
