# kali_rootfs/var/lib/dpkg/info/kali-defaults.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/postgresql@.service.d
/lib/systemd/system/postgresql@.service.d/kali_postgresql.conf
/etc
/etc/dconf
/etc/dconf/db
/etc/dconf/db/local.d
/etc/dconf/db/local.d/kali-tweaks
/etc/dconf/profile
/etc/dconf/profile/user
/etc/iceweasel
/etc/iceweasel/profile
/etc/iceweasel/pref
/etc/iceweasel/pref/kali.js
/etc/X11
/etc/X11/Xsession.d
/etc/X11/Xsession.d/52kali_fix_lang
/etc/X11/Xsession.d/52kali_noautomount
/etc/xdg
/etc/xdg/xfce4
/etc/xdg/xfce4/terminal
/etc/xdg/xfce4/terminal/terminalrc
/usr
/usr/share
/usr/share/kali-defaults
/usr/share/kali-defaults/.bashrc
/usr/share/kali-defaults/web
/usr/share/kali-defaults/web/js
/usr/share/kali-defaults/web/js/jquery.js
/usr/share/kali-defaults/web/js/bootstrap.min.js
/usr/share/kali-defaults/web/js/flexslider
/usr/share/kali-defaults/web/js/flexslider/jquery.flexslider.js
/usr/share/kali-defaults/web/js/flexslider/setting.js
/usr/share/kali-defaults/web/js/google-code-prettify
/usr/share/kali-defaults/web/js/google-code-prettify/prettify.js
/usr/share/kali-defaults/web/js/google-code-prettify/prettify.css
/usr/share/kali-defaults/web/js/jquery.flexslider.js
/usr/share/kali-defaults/web/js/jquery.fancybox-media.js
/usr/share/kali-defaults/web/js/portfolio
/usr/share/kali-defaults/web/js/portfolio/jquery.quicksand.js
/usr/share/kali-defaults/web/js/portfolio/setting.js
/usr/share/kali-defaults/web/js/validate.js
/usr/share/kali-defaults/web/js/quicksand
/usr/share/kali-defaults/web/js/quicksand/jquery.quicksand.js
/usr/share/kali-defaults/web/js/quicksand/setting.js
/usr/share/kali-defaults/web/js/jquery.fancybox.pack.js
/usr/share/kali-defaults/web/js/jquery.easing.1.3.js
/usr/share/kali-defaults/web/js/custom.js
/usr/share/kali-defaults/web/js/animate.js
/usr/share/kali-defaults/web/kali-site-logo.png
/usr/share/kali-defaults/web/skins
/usr/share/kali-defaults/web/skins/default.css
/usr/share/kali-defaults/web/edb-logo-gray-trans.png
/usr/share/kali-defaults/web/css
/usr/share/kali-defaults/web/css/fancybox
/usr/share/kali-defaults/web/css/fancybox/fancybox_sprite.png
/usr/share/kali-defaults/web/css/fancybox/fancybox_loading.gif
/usr/share/kali-defaults/web/css/fancybox/jquery.fancybox.css
/usr/share/kali-defaults/web/css/fancybox/blank.gif
/usr/share/kali-defaults/web/css/fancybox/fancybox_overlay.png
/usr/share/kali-defaults/web/css/custom-fonts.css
/usr/share/kali-defaults/web/css/style.css
/usr/share/kali-defaults/web/css/flexslider.css
/usr/share/kali-defaults/web/css/overwrite.css
/usr/share/kali-defaults/web/css/bootstrap.min.css
/usr/share/kali-defaults/web/css/font-awesome.css
/usr/share/kali-defaults/web/css/animate.css
/usr/share/kali-defaults/web/fonts
/usr/share/kali-defaults/web/fonts/FontAwesome.otf
/usr/share/kali-defaults/web/fonts/fontawesome-webfont.ttf
/usr/share/kali-defaults/web/fonts/fontawesome-webfont.eot
/usr/share/kali-defaults/web/fonts/fontawesome-webfont.woff
/usr/share/kali-defaults/web/fonts/fontawesome-webfont.svg
/usr/share/kali-defaults/web/fonts/customicon
/usr/share/kali-defaults/web/fonts/customicon/Icons.woff
/usr/share/kali-defaults/web/fonts/customicon/Icons.eot
/usr/share/kali-defaults/web/fonts/customicon/Icons.ttf
/usr/share/kali-defaults/web/fonts/customicon/Icons.svg
/usr/share/kali-defaults/web/kali-logo-gray-trans.png
/usr/share/kali-defaults/web/offsec-logo-gray-trans.png
/usr/share/kali-defaults/web/nethunter-logo-gray-trans.png
/usr/share/kali-defaults/web/img
/usr/share/kali-defaults/web/img/avatar.png
/usr/share/kali-defaults/web/img/works
/usr/share/kali-defaults/web/img/works/7.jpg
/usr/share/kali-defaults/web/img/works/8.jpg
/usr/share/kali-defaults/web/img/works/6.jpg
/usr/share/kali-defaults/web/img/works/3.jpg
/usr/share/kali-defaults/web/img/works/4.jpg
/usr/share/kali-defaults/web/img/works/5.jpg
/usr/share/kali-defaults/web/img/works/1.jpg
/usr/share/kali-defaults/web/img/works/2.jpg
/usr/share/kali-defaults/web/img/slides
/usr/share/kali-defaults/web/img/slides/3.jpg
/usr/share/kali-defaults/web/img/slides/1.jpg
/usr/share/kali-defaults/web/img/slides/2.jpg
/usr/share/kali-defaults/web/img/glyphicons-halflings.png
/usr/share/kali-defaults/web/img/glyphicons-halflings-white.png
/usr/share/kali-defaults/web/img/dummies
/usr/share/kali-defaults/web/img/dummies/dummy-1.jpg
/usr/share/kali-defaults/web/img/dummies/blog
/usr/share/kali-defaults/web/img/dummies/blog/img2.jpg
/usr/share/kali-defaults/web/img/dummies/blog/img1.jpg
/usr/share/kali-defaults/web/img/dummies/blog/65x65
/usr/share/kali-defaults/web/img/dummies/blog/65x65/thumb1.jpg
/usr/share/kali-defaults/web/img/dummies/blog/65x65/thumb3.jpg
/usr/share/kali-defaults/web/img/dummies/blog/65x65/thumb2.jpg
/usr/share/kali-defaults/web/img/dummies/blog/img3.jpg
/usr/share/kali-defaults/web/img/bg_direction_nav.png
/usr/share/kali-defaults/web/img/nivo-bullets.png
/usr/share/kali-defaults/web/homepage.html
/usr/share/kali-defaults/xsettings.xml
/usr/share/kali-defaults/bookmarks.html
/usr/share/kali-defaults/localstore.rdf
/usr/share/kali-defaults/postgresql_reduce_shared_buffers
/usr/share/glib-2.0
/usr/share/glib-2.0/schemas
/usr/share/glib-2.0/schemas/20_kali-defaults.gschema.override
/usr/share/doc
/usr/share/doc/kali-defaults
/usr/share/doc/kali-defaults/changelog.gz
/usr/share/doc/kali-defaults/copyright
