# kali_rootfs/var/lib/dpkg/info/network-manager.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/NetworkManager-wait-online.service
/lib/systemd/system/NetworkManager-dispatcher.service
/lib/systemd/system/NetworkManager.service
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/85-nm-unmanaged.rules
/lib/udev/rules.d/84-nm-drivers.rules
/usr
/usr/share
/usr/share/doc
/usr/share/doc/network-manager
/usr/share/doc/network-manager/changelog.Debian.gz
/usr/share/doc/network-manager/changelog.gz
/usr/share/doc/network-manager/README.gz
/usr/share/doc/network-manager/NEWS.Debian.gz
/usr/share/doc/network-manager/NEWS.gz
/usr/share/doc/network-manager/examples
/usr/share/doc/network-manager/examples/server.conf
/usr/share/doc/network-manager/copyright
/usr/share/doc/network-manager/README.Debian
/usr/share/doc/network-manager/AUTHORS
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/nmcli
/usr/share/polkit-1
/usr/share/polkit-1/rules.d
/usr/share/polkit-1/rules.d/60-network-manager.rules
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/org.freedesktop.NetworkManager.policy
/usr/share/dbus-1
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.nm_dispatcher.service
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/NetworkManager.8.gz
/usr/share/man/man5
/usr/share/man/man5/nmcli-examples.5.gz
/usr/share/man/man5/nm-system-settings.conf.5.gz
/usr/share/man/man5/nm-settings.5.gz
/usr/share/man/man5/nm-settings-keyfile.5.gz
/usr/share/man/man5/NetworkManager.conf.5.gz
/usr/share/man/man1
/usr/share/man/man1/nmtui.1.gz
/usr/share/man/man1/nmtui-hostname.1.gz
/usr/share/man/man1/nmtui-edit.1.gz
/usr/share/man/man1/nmtui-connect.1.gz
/usr/share/man/man1/nmcli.1.gz
/usr/share/man/man1/nm-online.1.gz
/usr/share/locale
/usr/share/locale/zh_TW
/usr/share/locale/zh_TW/LC_MESSAGES
/usr/share/locale/zh_TW/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/zh_HK
/usr/share/locale/zh_HK/LC_MESSAGES
/usr/share/locale/zh_HK/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/zh_CN
/usr/share/locale/zh_CN/LC_MESSAGES
/usr/share/locale/zh_CN/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/wa
/usr/share/locale/wa/LC_MESSAGES
/usr/share/locale/wa/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/vi
/usr/share/locale/vi/LC_MESSAGES
/usr/share/locale/vi/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/uk
/usr/share/locale/uk/LC_MESSAGES
/usr/share/locale/uk/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/tr
/usr/share/locale/tr/LC_MESSAGES
/usr/share/locale/tr/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/th
/usr/share/locale/th/LC_MESSAGES
/usr/share/locale/th/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/te
/usr/share/locale/te/LC_MESSAGES
/usr/share/locale/te/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/ta
/usr/share/locale/ta/LC_MESSAGES
/usr/share/locale/ta/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/sv
/usr/share/locale/sv/LC_MESSAGES
/usr/share/locale/sv/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/sr@latin
/usr/share/locale/sr@latin/LC_MESSAGES
/usr/share/locale/sr@latin/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/sr
/usr/share/locale/sr/LC_MESSAGES
/usr/share/locale/sr/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/sq
/usr/share/locale/sq/LC_MESSAGES
/usr/share/locale/sq/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/sl
/usr/share/locale/sl/LC_MESSAGES
/usr/share/locale/sl/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/sk
/usr/share/locale/sk/LC_MESSAGES
/usr/share/locale/sk/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/rw
/usr/share/locale/rw/LC_MESSAGES
/usr/share/locale/rw/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/ru
/usr/share/locale/ru/LC_MESSAGES
/usr/share/locale/ru/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/pt_BR
/usr/share/locale/pt_BR/LC_MESSAGES
/usr/share/locale/pt_BR/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/pt
/usr/share/locale/pt/LC_MESSAGES
/usr/share/locale/pt/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/pl
/usr/share/locale/pl/LC_MESSAGES
/usr/share/locale/pl/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/pa
/usr/share/locale/pa/LC_MESSAGES
/usr/share/locale/pa/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/or
/usr/share/locale/or/LC_MESSAGES
/usr/share/locale/or/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/oc
/usr/share/locale/oc/LC_MESSAGES
/usr/share/locale/oc/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/nl
/usr/share/locale/nl/LC_MESSAGES
/usr/share/locale/nl/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/ne
/usr/share/locale/ne/LC_MESSAGES
/usr/share/locale/ne/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/nb
/usr/share/locale/nb/LC_MESSAGES
/usr/share/locale/nb/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/mr
/usr/share/locale/mr/LC_MESSAGES
/usr/share/locale/mr/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/ml
/usr/share/locale/ml/LC_MESSAGES
/usr/share/locale/ml/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/mk
/usr/share/locale/mk/LC_MESSAGES
/usr/share/locale/mk/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/lv
/usr/share/locale/lv/LC_MESSAGES
/usr/share/locale/lv/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/lt
/usr/share/locale/lt/LC_MESSAGES
/usr/share/locale/lt/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/ku
/usr/share/locale/ku/LC_MESSAGES
/usr/share/locale/ku/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/ko
/usr/share/locale/ko/LC_MESSAGES
/usr/share/locale/ko/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/kn
/usr/share/locale/kn/LC_MESSAGES
/usr/share/locale/kn/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/ka
/usr/share/locale/ka/LC_MESSAGES
/usr/share/locale/ka/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/ja
/usr/share/locale/ja/LC_MESSAGES
/usr/share/locale/ja/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/it
/usr/share/locale/it/LC_MESSAGES
/usr/share/locale/it/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/id
/usr/share/locale/id/LC_MESSAGES
/usr/share/locale/id/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/hu
/usr/share/locale/hu/LC_MESSAGES
/usr/share/locale/hu/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/hr
/usr/share/locale/hr/LC_MESSAGES
/usr/share/locale/hr/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/hi
/usr/share/locale/hi/LC_MESSAGES
/usr/share/locale/hi/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/he
/usr/share/locale/he/LC_MESSAGES
/usr/share/locale/he/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/gu
/usr/share/locale/gu/LC_MESSAGES
/usr/share/locale/gu/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/gl
/usr/share/locale/gl/LC_MESSAGES
/usr/share/locale/gl/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/gd
/usr/share/locale/gd/LC_MESSAGES
/usr/share/locale/gd/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/fr
/usr/share/locale/fr/LC_MESSAGES
/usr/share/locale/fr/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/fi
/usr/share/locale/fi/LC_MESSAGES
/usr/share/locale/fi/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/eu
/usr/share/locale/eu/LC_MESSAGES
/usr/share/locale/eu/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/et
/usr/share/locale/et/LC_MESSAGES
/usr/share/locale/et/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/es
/usr/share/locale/es/LC_MESSAGES
/usr/share/locale/es/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/eo
/usr/share/locale/eo/LC_MESSAGES
/usr/share/locale/eo/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/en_GB
/usr/share/locale/en_GB/LC_MESSAGES
/usr/share/locale/en_GB/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/en_CA
/usr/share/locale/en_CA/LC_MESSAGES
/usr/share/locale/en_CA/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/el
/usr/share/locale/el/LC_MESSAGES
/usr/share/locale/el/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/dz
/usr/share/locale/dz/LC_MESSAGES
/usr/share/locale/dz/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/de
/usr/share/locale/de/LC_MESSAGES
/usr/share/locale/de/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/da
/usr/share/locale/da/LC_MESSAGES
/usr/share/locale/da/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/cs
/usr/share/locale/cs/LC_MESSAGES
/usr/share/locale/cs/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/ca
/usr/share/locale/ca/LC_MESSAGES
/usr/share/locale/ca/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/bs
/usr/share/locale/bs/LC_MESSAGES
/usr/share/locale/bs/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/bn_IN
/usr/share/locale/bn_IN/LC_MESSAGES
/usr/share/locale/bn_IN/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/bg
/usr/share/locale/bg/LC_MESSAGES
/usr/share/locale/bg/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/be@latin
/usr/share/locale/be@latin/LC_MESSAGES
/usr/share/locale/be@latin/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/as
/usr/share/locale/as/LC_MESSAGES
/usr/share/locale/as/LC_MESSAGES/NetworkManager.mo
/usr/share/locale/ar
/usr/share/locale/ar/LC_MESSAGES
/usr/share/locale/ar/LC_MESSAGES/NetworkManager.mo
/usr/lib
/usr/lib/pppd
/usr/lib/pppd/2.4.7
/usr/lib/pppd/2.4.7/nm-pppd-plugin.so
/usr/lib/arm-linux-gnueabihf
/usr/lib/arm-linux-gnueabihf/NetworkManager
/usr/lib/arm-linux-gnueabihf/NetworkManager/libnm-settings-plugin-ifupdown.so
/usr/lib/arm-linux-gnueabihf/NetworkManager/libnm-settings-plugin-ibft.so
/usr/lib/arm-linux-gnueabihf/NetworkManager/libnm-device-plugin-adsl.so
/usr/lib/arm-linux-gnueabihf/NetworkManager/libnm-device-plugin-bluetooth.so
/usr/lib/arm-linux-gnueabihf/NetworkManager/libnm-device-plugin-team.so
/usr/lib/arm-linux-gnueabihf/NetworkManager/libnm-device-plugin-wifi.so
/usr/lib/arm-linux-gnueabihf/NetworkManager/libnm-device-plugin-wwan.so
/usr/lib/arm-linux-gnueabihf/NetworkManager/libnm-wwan.so
/usr/lib/NetworkManager
/usr/lib/NetworkManager/nm-dhcp-helper
/usr/lib/NetworkManager/nm-iface-helper
/usr/lib/NetworkManager/nm-dispatcher
/usr/lib/NetworkManager/ifblacklist_migrate.sh
/usr/bin
/usr/bin/nm-online
/usr/bin/nmcli
/usr/bin/nmtui
/usr/sbin
/usr/sbin/NetworkManager
/etc
/etc/init.d
/etc/init.d/network-manager
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/nm-dispatcher.conf
/etc/dbus-1/system.d/org.freedesktop.NetworkManager.conf
/etc/NetworkManager
/etc/NetworkManager/NetworkManager.conf
/etc/NetworkManager/dispatcher.d
/etc/NetworkManager/dispatcher.d/01ifupdown
/etc/NetworkManager/dispatcher.d/pre-up.d
/etc/NetworkManager/dispatcher.d/pre-down.d
/etc/NetworkManager/dispatcher.d/no-wait.d
/etc/NetworkManager/VPN
/etc/NetworkManager/system-connections
/etc/NetworkManager/dnsmasq-shared.d
/etc/NetworkManager/dnsmasq.d
/etc/NetworkManager/conf.d
/var
/var/lib
/var/lib/polkit-1
/var/lib/polkit-1/localauthority
/var/lib/polkit-1/localauthority/10-vendor.d
/var/lib/polkit-1/localauthority/10-vendor.d/org.freedesktop.NetworkManager.pkla
/var/lib/NetworkManager
/lib/systemd/system/network-manager.service
/usr/bin/nmtui-hostname
/usr/bin/nmtui-edit
/usr/bin/nmtui-connect
