# kali_rootfs/var/lib/dpkg/info/paxctld.conffiles
/etc/init.d/paxctld
/etc/init/paxctld.conf
/etc/paxctld.conf
