# kali_rootfs/var/lib/dpkg/info/paxctld.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/paxctld.service
/usr
/usr/share
/usr/share/doc
/usr/share/doc/paxctld
/usr/share/doc/paxctld/changelog.Debian.gz
/usr/share/doc/paxctld/copyright
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/paxctld.8.gz
/sbin
/sbin/paxctld
/etc
/etc/init.d
/etc/init.d/paxctld
/etc/init
/etc/init/paxctld.conf
/etc/paxctld.conf
