# kali_rootfs/lib/systemd/system/postgresql@.service.d/kali_postgresql.conf
[Service]
ExecStartPre=/usr/share/kali-defaults/postgresql_reduce_shared_buffers %I
